# templates/slack-mastodon.timer.erb
[Unit]
Description=Schedule mirroring mastodon to slack

[Timer]
OnBootSec=<%= @bootdelay %>
OnUnitActiveSec=<%= @frequency %>

[Install]
WantedBy=multi-user.target
