# files/slack-mastodon.service
[Unit]
Description=Mirror Mastodon to Slack

[Service]
Type=simple
User=slackmastodon
ExecStart=/usr/local/bin/slack-mastodon update -v
