# files/slack-mastodon.service
[Unit]
Description=Mirror Mastodon to Slack

[Service]
Type=simple
ExecStart=/usr/local/bin/slack-mastodon update -v
